# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f

[program: streamlit]
command=poetry run streamlit run app/src/inference_pipeline/frontend/main.py --server.port 8501